# build/etc/supervisor.d/nagios.ini
[program:nagios]
command=nagios /etc/nagios/nagios.cfg
autostart=true
autorestart=true
user=nagios
priority=2
startretries=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log
stdout_logfile_maxbytes=0

